# neuroling.service
[Unit]
Description=Launch the speech-to-text  

[Service]
Type=simple
ExecStart=/home/myserver/bin/neuroling.sh

